# deployment/roles/ckanext-harvest/templates/etc/cron.d/ckan-harvest
*/2 * * * * root /usr/lib/ckan/default/bin/paster --plugin=ckanext-harvest harvester run -c /etc/ckan/default/production.ini >> /var/log/harvester_run.log 2>&1
00 05 * * * root /usr/lib/ckan/default/bin/paster --plugin=ckanext-glasgow db_clean harvest -c /etc/ckan/default/production.ini >> /var/log/harvester_run.log 2>&1
